# aliran-tunai.service
[Unit]
Description=AliranTunai API Server Service
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/aliran-tunai/current
Environment=PATH=/opt/aliran-tunai/current/venv/bin
ExecStart=/opt/aliran-tunai/current/venv/bin/python api_server.py
Restart=always
RestartSec=5

# Environment variables
EnvironmentFile=/opt/aliran-tunai/current/.env

# Security settings - Comprehensive hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/aliran-tunai
ReadWritePaths=/var/log/aliran-tunai

# Network security
PrivateNetwork=no
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=10.0.0.0/8
IPAddressAllow=172.16.0.0/12
IPAddressAllow=192.168.0.0/16

# Filesystem security
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ReadOnlyPaths=/
ReadWritePaths=/opt/aliran-tunai
ReadWritePaths=/tmp
ReadWritePaths=/var/tmp
ReadWritePaths=/var/log/aliran-tunai

# Process security
PrivateDevices=yes
ProtectHostname=yes
ProtectClock=yes
RemoveIPC=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes

# System calls security
SystemCallFilter=@system-service
SystemCallFilter=~@debug @mount @cpu-emulation @obsolete @privileged @reboot @swap @raw-io
SystemCallErrorNumber=EPERM

# Capabilities (drop all capabilities)
CapabilityBoundingSet=
AmbientCapabilities=

# Memory protection
MemoryDenyWriteExecute=yes
ProtectProc=invisible
ProcSubset=pid

# Resource limits
MemoryMax=512M
CPUQuota=80%
TasksMax=100
LimitNOFILE=1024
LimitNPROC=50

# Timeouts
TimeoutStartSec=30s
TimeoutStopSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aliran-tunai

[Install]
WantedBy=multi-user.target
