# aliran-tunai.service
[Unit]
Description=AliranTunai Telegram Bot
After=network.target

[Service]
Type=simple
User=megat-aliran
Group=megat-aliran
WorkingDirectory=/home/megat-aliran/aliran-tunai
Environment=PATH=/home/megat-aliran/aliran-tunai/aliran/bin
ExecStart=/home/megat-aliran/aliran-tunai/aliran/bin/python main.py
Restart=always
RestartSec=3

# Environment variables
EnvironmentFile=/home/megat-aliran/aliran-tunai/.env

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/megat-aliran/aliran-tunai

[Install]
WantedBy=multi-user.target
