# deploy/aliran-tunai.service
[Unit]
Description=AliranTunai Telegram Bot Service
After=network.target mongod.service
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/aliran-tunai/current
Environment=PATH=/opt/aliran-tunai/current/venv/bin
ExecStart=/opt/aliran-tunai/current/venv/bin/python main.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aliran-tunai

# Security
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/opt/aliran-tunai/current
ReadWritePaths=/var/log/aliran-tunai
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65536
MemoryLimit=1G
CPUQuota=50%

[Install]
WantedBy=multi-user.target